# supervisor.conf
[program:yams]
command = /home/pi/yams/env/bin/python -m yams
directory = /home/pi/yams
startsecs = 5
redirect_stderr = True
autostart = True
autorestart = True
